# bin/ego_browse_passenger
#!/bin/bash 


#	Copy this script to /etc/init.d/
#
#	sudo cp $0 /etc/init.d/
#	sudo update-rc.d ego_browse_passenger defaults


### BEGIN INIT INFO 
# Provides:					EGO Browse passenger
# Default-Start:		 2 3 4 5 
# Default-Stop:			0 1 6 
# Required-Start:		$remote_fs $syslog
# Required-Stop:		 $remote_fs $syslog
# X-Interactive:		 true 
# Short-Description: Start/stop ego_browse web site 
### END INIT INFO 

#PATH=/sbin:/usr/sbin:/bin:/usr/bin 
PATH=/home/ruby/.local/bin:/home/ruby/.gem/ruby/2.3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
export GEM_HOME=/home/ruby/.gem/ruby/2.3
PASSENGER=/home/ruby/.gem/ruby/2.3/bin/passenger
ADDRESS=127.0.0.1 
PORT=3000 
RAILS_ENV=production
ENVIRONMENT=production 
APP=/var/www/ego_browse
#USER=ruby


. /lib/lsb/init-functions 


case "$1" in 
	start) 
		echo "Starting EGO Web passenger" 
		#		iptables --policy INPUT ACCEPT
		#	Be more specific
		iptables -A INPUT -p tcp --dport 3000 -m conntrack --ctstate NEW,ESTABLISHED -j ACCEPT
		cd $APP
		echo "$PASSENGER start --port $PORT --environment $ENVIRONMENT --daemonize"
		$PASSENGER start --port $PORT --environment $ENVIRONMENT --daemonize
		;; 
	stop) 
		echo "Stopping EGO Web passenger" 
		cd $APP
		echo "$PASSENGER stop -p $PORT"
		$PASSENGER stop -p $PORT
		;; 
	status)
		#	From	/lib/lsb/init-functions 
		status_of_proc -p $APP/tmp/pids/passenger.$PORT.pid "" "passenger" && exit 0 || exit $?
		;;
	*) 
		echo "Usage: $0 start|stop" >&2 
		exit 3 
		;; 
esac
